# dsa-nagios-nrpe-config/debian/cron.daily
#!/bin/sh

[ -x /usr/share/dsa/weak-ssh-keys-check ] && /usr/share/dsa/weak-ssh-keys-check -s /var/cache/dsa/nagios/weak-ssh-keys
